# scripts/deploy/personalwebsite.service
[Unit]
Description=PersonalWebsite FastAPI service
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=PersonalWebsite/backend
Environment=PORT=8000
Environment=APP_DOMAIN=https://nathanieljshepherd.com
Environment=ALLOWED_ORIGINS=https://nathanieljshepherd.com,https://www.nathanieljshepherd.com
ExecStart=PersonalWebsite/backend/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target